# crontab
59 23 * * * /usr/local/bin/python /app/brain.py >> /var/log/cron.log 2>&1